# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:xinetd]
command=service xinetd restart
autorestart=false
